# daemon/deployment/tau-daemon.service
[Unit]
Description=Tau Lighting Control Daemon
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=tau
Group=tau
WorkingDirectory=/opt/tau-daemon/daemon
EnvironmentFile=/opt/tau-daemon/daemon/.env
Environment="PYTHONPATH=/opt/tau-daemon/daemon/src"

ExecStart=/opt/tau-daemon/daemon/.venv/bin/python -m tau.main
Restart=on-failure
RestartSec=10s

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/log/tau /opt/tau-daemon

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tau-daemon

[Install]
WantedBy=multi-user.target
