# daemon/deployment/tau-daemon.service
[Unit]
Description=Tau Lighting Control Daemon
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=tau
Group=tau
WorkingDirectory=/opt/tau-daemon
Environment="PYTHONPATH=/opt/tau-daemon/src"
Environment="DATABASE_URL=postgresql://tau_daemon:tau_password@localhost/tau_lighting"
Environment="DAEMON_PORT=8000"
Environment="LOG_LEVEL=INFO"
Environment="CONTROL_LOOP_HZ=30"

ExecStart=/opt/tau-daemon/.venv/bin/python -m tau.main
Restart=on-failure
RestartSec=10s

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/tau /opt/tau-daemon

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tau-daemon

[Install]
WantedBy=multi-user.target
